# gongxunsai.service
[Unit]
Description=Gongxun Xiao Car
After=network.target

[Service]
User=pi
Group=pi
ExecStart=/usr/bin/python3 /home/pi/gongxunsai/app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
